# i3.target

[Unit]
Description=i3 desktop
Requires=dbus.socket
AllowIsolate=true

