# i3wm.service

[Unit]
Description=i3 Window Manager
After=dbus.socket

[Service]
ExecStart=/usr/bin/i3

